# hand-holding-detector.service
[Unit]
Description=Hand Holding Detector ROS2 Node
After=network.target ros2-daemon.service
Wants=network.target

[Service]
Type=simple
User=sunrise
Group=sunrise
WorkingDirectory=/home/sunrise/Desktop/bloom
Environment="HOME=/home/sunrise"
Environment="USER=sunrise"
ExecStart=/home/sunrise/Desktop/bloom/start_hand_holding_detector.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits (optional, adjust as needed)
# LimitCPU=2
# LimitMEM=512M

[Install]
WantedBy=multi-user.target


